# supervisord.conf
; append this section to /etc/supervisord.conf
[program:subfusc]
user=bgraham
command=python /var/www/subfusc/Subfusc.py